# services/node/run
#!/bin/bash

docker-compose run --user $(id -u):$(id -g) --rm node sh -c "$@"
